# roles/idazmaia/files/etc/xdg/autostart/proiektorea.desktop
[Desktop Entry]
Name=Proiektorea konfiguratu
Exec=/usr/local/bin/proiektorea
Icon=video-display
