# roles/ikastola/files/etc/systemd/system/itzali.service
[Unit]
Description=Sistema itzali

[Service]
Type=oneshot
ExecStart=/usr/sbin/poweroff
