# konfigurazioak/etc/xdg/autostart/updatehosts.desktop
[Desktop Entry]
Name=Hosts eguneratu
Exec=/usr/local/bin/update_hosts
Icon=applications-internet
